# etc/profile.d/welcome.sh
sh /usr/bin/deersay.sh "Welcome to DeerPi for Rpi0"
